# db/db-timer/jand-db-timer.service
[Unit]
Description=JAND database timer.

[Service]
ExecStart=/srv/jand-db-timer/jand-db-timer.php
User=jand
Group=jand
